# server/utils/docker/services/sensorcfg-creation-worker/run
#!/bin/bash
python /opt/HoneySens/beanstalk/honeysens-sensorcfg-creation-worker.py /opt/HoneySens/data/config.cfg
