# server/utils/docker/services/update-worker/run
#!/bin/bash
python /opt/HoneySens/beanstalk/honeysens-update-worker.py /opt/HoneySens/data/config.cfg
