# sensor/platform/bbb/debian/postinst
#!/bin/sh
# postinst script for honeysens-sensor-platform-bbb
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        # Disable sshd on boot and use our own service that triggers after the USB Ethernet gadget is available
        systemctl disable ssh
        systemctl enable honeysens-sensor-platform-bbb-sshd
	# Run sshd only on the USB ethernet device
	sed -i "s/#ListenAddress 0.0.0.0/ListenAddress 192.168.7.2/g" /etc/ssh/sshd_config
	# Disable systemd's NTP, because it interferes with HTTPS time sync
	systemctl disable systemd-timesyncd
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
